# deploy/systemd/dongchedi-update.timer
[Unit]
Description=Run Dongchedi China catalog update daily (same time as Korea)

[Timer]
# Полночь Asia/Yekaterinburg — как rideauto-auto-update.timer
OnCalendar=*-*-* 00:00:00 Asia/Yekaterinburg
Persistent=true
Unit=dongchedi-update.service

[Install]
WantedBy=timers.target
