# deploy/systemd/dongchedi-update.service
[Unit]
Description=Prod Encar Dongchedi → Postgres or encar_china.db (dongchedi_scraper.yaml storage)
After=network.target
# Первый раз: полный прогон вручную (см. deploy/DEPLOY.md). Далее — ежедневно в полночь по Екатеринбургу (dongchedi-update.timer).
# Права: как у encar-update — www-data пишет encar_china.db

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/prod-encar/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/prod-encar-scrapers
ExecStart=/opt/prod-encar/.venv/bin/python -m dongchedi.scraper --config /opt/prod-encar/backend/dongchedi_scraper.yaml --db /opt/prod-encar/encar_china.db
# При storage.backend=postgres задайте DATABASE_URL в окружении (EnvironmentFile) или в YAML.
TimeoutStartSec=24h
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/prod-encar

[Install]
WantedBy=multi-user.target
