# deploy/systemd/encar-update.timer
[Unit]
Description=Run Prod Encar update daily

[Timer]
# Полночь Asia/Yekaterinburg — encar_daily_update / auto_update (Postgres; опционально postgres_catalog_sync)
OnCalendar=*-*-* 00:00:00 Asia/Yekaterinburg
Persistent=true
Unit=encar-update.service

[Install]
WantedBy=timers.target
