# deploy/systemd/encar-update.timer
[Unit]
Description=Run Prod Encar update daily

[Timer]
# 12:00 Asia/Yekaterinburg (ежедневное обновление: новые/проданные → цены → фронт)
OnCalendar=*-*-* 12:00:00 Asia/Yekaterinburg
Persistent=true
Unit=encar-update.service

[Install]
WantedBy=timers.target
