# deploy/systemd/encar-update.service
[Unit]
Description=Prod Encar daily update job
After=network.target
# Права: User=www-data — logs/, доступ к Postgres (сеть/credentials), опционально локальные *.db после миграции
# sudo bash /opt/prod-encar/deploy/scripts/ensure_scraper_runtime_permissions.sh

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/prod-encar
Environment=PYTHONUNBUFFERED=1
# DATABASE_URL для encar_scraper (storage.postgres) и Checkpoint в encar_daily_update.
EnvironmentFile=-/etc/default/prod-encar-scrapers
# Не auto_update.py: он тянет legacy PostgreSQLDatabase (inner_id), несовместим с продовой таблицей cars (car_id).
ExecStart=/opt/prod-encar/.venv/bin/python /opt/prod-encar/backend/encar_daily_update.py --once --config /opt/prod-encar/scraper_config.yaml
TimeoutStartSec=24h
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/prod-encar
