# deploy/systemd/encar-subscriptions-notify.service
[Unit]
Description=Prod Encar subscription notifications runner
After=network.target encar-api.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/prod-encar
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/prod-encar
ExecStart=/opt/prod-encar/.venv/bin/python /opt/prod-encar/backend/scripts/run_subscription_notifications.py --api http://127.0.0.1:8080 --admin-key ${SUBSCRIPTIONS_ADMIN_KEY}
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/prod-encar
