# deploy/systemd/prod-encar-auto-update.service
[Unit]
Description=Prod Encar Auto Update (PostgreSQL)
After=network.target

[Service]
Type=oneshot
User=prod-encar
Group=prod-encar
WorkingDirectory=/opt/prod-encar
Environment=PYTHONUNBUFFERED=1
# Ускорить ночной цикл: в /etc/default/prod-encar можно задать SKIP_LEARN_ENGINE_MAP=1.
EnvironmentFile=-/etc/default/prod-encar
# Совпадает с encar-update.service: только encar_daily_update (каталог в public.cars), без legacy PostgreSQLDatabase.
ExecStart=/opt/prod-encar/.venv/bin/python /opt/prod-encar/backend/encar_daily_update.py --once --config /opt/prod-encar/scraper_config.yaml
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/prod-encar
LockPersonality=true
MemoryDenyWriteExecute=true
