# deploy/systemd/prod-dongchedi-update.service
[Unit]
Description=RideAuto Dongchedi catalog → PostgreSQL or encar_china.db (see dongchedi_scraper.yaml storage)
After=network.target
# Аналог dongchedi-update.service для стека rideauto (см. rideauto-auto-update.service).

[Service]
Type=oneshot
User=rideauto
Group=rideauto
WorkingDirectory=/opt/rideauto/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/rideauto
ExecStart=/usr/bin/python3 -m dongchedi.scraper --config /opt/rideauto/backend/dongchedi_scraper.yaml --db /opt/rideauto/encar_china.db
# DATABASE_URL в EnvironmentFile (rideauto) для storage.backend=postgres; --db задаёт имя checkpoint (*.scraper.checkpoint.json).
TimeoutStartSec=24h
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/rideauto
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
