# deploy/systemd/prod-dongchedi-update.timer
[Unit]
Description=Run Dongchedi update daily (rideauto, same time as Korea)

[Timer]
# На час позже Кореи (rideauto-auto-update), чтобы не бить API и Postgres одновременно.
OnCalendar=*-*-* 01:00:00 Asia/Yekaterinburg
Persistent=true
Unit=prod-dongchedi-update.service

[Install]
WantedBy=timers.target
