# deploy/systemd/prod-encar-meilisearch-sync.service
[Unit]
Description=Prod Encar — PostgreSQL → Meilisearch (ночная выгрузка каталога на фронт)
After=network.target

[Service]
Type=oneshot
User=prod-encar
Group=prod-encar
WorkingDirectory=/opt/prod-encar
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/prod-encar
ExecStart=/bin/bash /opt/prod-encar/deploy/scripts/run_meilisearch_sync_host.sh
TimeoutStartSec=6h
NoNewPrivileges=true
PrivateTmp=true
