# deploy/systemd/prod-encar-auto-update.service
[Unit]
Description=Prod Encar Auto Update (PostgreSQL)
After=network.target

[Service]
Type=oneshot
User=prod-encar
Group=prod-encar
WorkingDirectory=/opt/prod-encar
Environment=PYTHONUNBUFFERED=1
# Ускорить ночной цикл: в /etc/default/prod-encar можно задать SKIP_LEARN_ENGINE_MAP=1.
EnvironmentFile=-/etc/default/prod-encar
# Единственный ночной пайплайн Encar для rideauto: discover + sold + encar_scraper --only-pending (см. encar_daily_update.py).
ExecStart=/opt/prod-encar/.venv/bin/python /opt/prod-encar/backend/encar_daily_update.py --once --config /opt/prod-encar/scraper_config.yaml
# Без этого systemd по умолчанию (часто ~90s) убивает длинный scraper до конца цикла.
TimeoutStartSec=24h
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/prod-encar
LockPersonality=true
MemoryDenyWriteExecute=true
