# deploy/systemd/rideauto-hp-catalog-fill.service
[Unit]
Description=RideAuto HP catalog LLM filler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rideauto
Group=rideauto
WorkingDirectory=/opt/rideauto
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/rideauto
ExecStart=/opt/rideauto/.venv/bin/python /opt/rideauto/backend/scripts/fill_hp_catalog_deepseek.py --db /opt/rideauto/data/hp_catalog.db --provider auto --model deepseek-chat --openai-model gpt-4o-mini --retry-errors --continuous --sleep-sec 0.5 --idle-sleep-sec 20 --max-attempts 8
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/rideauto
LockPersonality=true
MemoryDenyWriteExecute=true

