# deploy/systemd/rideauto-sitemap-gen.timer
[Unit]
Description=Daily disk sitemap generation (rideauto)

[Timer]
OnCalendar=*-*-* 04:30:00 Asia/Yekaterinburg
Persistent=true
Unit=rideauto-sitemap-gen.service

[Install]
WantedBy=timers.target
