# deploy/systemd/rideauto-subscriptions-notify.timer
[Unit]
Description=Run RideAuto subscription notifications every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=rideauto-subscriptions-notify.service

[Install]
WantedBy=timers.target
