# deploy/systemd/rideauto-auto-update.timer
[Unit]
Description=Run RideAuto auto update daily

[Timer]
# Полночь по Екатеринбургу: новые объявления, выборочная проверка проданных, докачка деталей, export + цены → фронт
OnCalendar=*-*-* 00:00:00 Asia/Yekaterinburg
Persistent=true
Unit=rideauto-auto-update.service

[Install]
WantedBy=timers.target
