# deploy/encar-api.service.example
# Скопируйте в /etc/systemd/system/encar-api.service и подставьте пути, пользователя и порт.
# Переменные окружения (DSN Postgres, Meilisearch, Redis): см. /etc/default/prod-encar или Environment= ниже.
# После правок: sudo systemctl daemon-reload && sudo systemctl enable --now encar-api

[Unit]
Description=Prod Encar catalog API (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/prod-encar/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/prod-encar
# Environment=WRA_PG_DSN=postgresql://...
ExecStart=/var/www/prod-encar/backend/.venv/bin/python -m uvicorn fastapi_app.main:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
